# scripts/esp32-listener.service
[Unit]
Description=ESP32 NeoPixel Contact Closure Listener
After=network.target

[Service]
Type=simple
User=fpp
Group=fpp
WorkingDirectory=/home/fpp/media/plugins/contact_clsr_fpp
ExecStart=/home/fpp/media/plugins/contact_clsr_fpp/scripts/esp32_listener.sh
Restart=always
RestartSec=5
StandardOutput=append:/home/fpp/media/logs/contact_clsr_fpp.log
StandardError=append:/home/fpp/media/logs/contact_clsr_fpp.log

[Install]
WantedBy=multi-user.target

